# php/shared/etc/confd/templates/supervisor/dhparam_regenerate.conf.tmpl
[program:dhparam-regenerate]
command = supervisorctl --user supervisor --password supervisor start dhparam-regenerate
startsecs = 0
autostart = no
autorestart = false