# php/shared/etc/confd/templates/supervisor/dhparam_regenerate.conf.tmpl
[program:dhparam-regenerate]
command = container --no-debug dhparam_regenerate
startsecs = 0
autostart = no
autorestart = false