# php-apache/7.0/etc/supervisor/conf.d/apache.conf
[program:apache]
priority = 2
command = /etc/apache2/run.sh
autostart = true
user = root
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
